# etc/supervisord.d/consul.ini
[program:consul]
command=/opt/qnib/bin/start_consul.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor_%(program_name)s.log
autorestart = false
autostart = true
startsecs=2
exitcodes=0
