# etc/supervisord.d/consul.ini
[program:consul]
command=/opt/qnib/bin/start_consul.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_events_enabled = false
stderr_events_enabled = false
autorestart = false
autostart = true
startsecs=2
exitcodes=0
