# etc/supervisord.d/consul.ini
[program:consul]
command=/opt/qnib/consul/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_events_enabled = false
stderr_events_enabled = false
autorestart = false
autostart = true
startsecs=5
startretries=50
exitcodes=0
